# twito/conf/gunicorn_supervisor.conf
[program:celery]
command = /home/worker/twito/twito/conf/django_celery_start.bash
user = worker
stfout_logfile = /home/worker/twito/twito/logs/django_celery.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[program:twito]
command = /home/worker/twito/twito/conf/gunicorn_start.bash                    ; Command to start app
user = worker                                                                  ; User to run as
stdout_logfile = /home/worker/twito/twito/logs/gunicorn_supervisor.log         ; Where to write log messages
redirect_stderr = true                                                         ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                                ; Set UTF-8 as default encoding
